# build/config/upstart/mount_ebslvm.conf
# upstart file to mount volume

start on filesystem

emits ebslvm-mounted

script
  /usr/sbin/ebslvm vgebs lvebs /mnt
  chmod g+w /mnt
  chgrp edge /mnt
  initctl emit ebslvm-mounted
end script
