# build/config/upstart/spade_edge.conf
# Upstart script at /etc/init/spade_edge.conf

start on ebslvm-mounted and stopped cloud-final
# supposedly this will cause the system to hang shutdown on
# upstart killing the edge: http://goo.gl/lJQVzH
stop on starting rc RUNLEVEL=[016]

emits edge-running

respawn
respawn limit 10 5

limit nofile 4096 4096

kill timeout 300
kill signal SIGINT

script
  # use su <user> because that makes /etc/environment variables available to subprocess
  exec su root -c /opt/science/spade_edge/bin/run_edge.sh
  initctl emit edge-running
end script

