# build/config/mount_ebslvm.conf
# upstart file to mount volume

start on filesystem

emits ebslvm-mounted

script
  /usr/sbin/ebslvm vgebs lvebs /mnt
  initctl emit ebslvm-mounted
end script
